# supervisor-app.conf
[program:node-server]
command = /usr/local/bin/npm start

[program:nginx-app]
command = /usr/sbin/nginx
